# supervisord.conf
# ─────────────────────────────────────────────
# BioOps Twin — Process Orchestration
# ─────────────────────────────────────────────
# Manages two concurrent processes for the Hugging Face Space:
#   1. mock_proxy  — Veea Lobster Trap security proxy (port 8080)
#   2. gradio      — Main Gradio dashboard (port 7860)
# ─────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:mock_proxy]
command=python mock_veea_proxy.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:gradio]
command=python main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
